# faucet/server/wintg-faucet.service
[Unit]
Description=WINTG Faucet API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wintg-faucet
EnvironmentFile=/opt/wintg-faucet/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
User=wintg
Group=wintg
StandardOutput=append:/var/log/wintg-faucet.log
StandardError=append:/var/log/wintg-faucet.log

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log /opt/wintg-faucet

[Install]
WantedBy=multi-user.target
